# app/config/crontab/production
# parser price skins
0   0   *   *   *   /usr/bin/flock -n /tmp/parser_market /usr/bin/php /var/www/cgrays/current/bin/console processing:parser_market -vvv 730 --env=prod --no-debug >> /tmp/parser_market.log 2>&1
*   *   *   *   *   /usr/bin/flock -n /tmp/parser_skins /usr/bin/php /var/www/cgrays/current/bin/console rabbitmq:consumer parser_skins -w -m 10 -vvv --env=prod --no-debug >> /tmp/parser_skins.log 2>&1
*   *   *   *   *   /usr/bin/flock -n /tmp/parser_price /usr/bin/php /var/www/cgrays/current/bin/console rabbitmq:consumer parser_price -w -m 1000 -vvv --env=prod --no-debug >> /tmp/parser_price.log 2>&1
